# docker/websockets.conf
[program:websockets]
command=/usr/bin/php /var/www/artisan websocket$
numprocs=1
autostart=true
autorestart=true
user=root
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log
